# sock2port/sock2port@.socket
[Unit]
Description=Forward port %i traffic to Unix socket

[Socket]
ListenStream=%i

[Install]
WantedBy=sockets.target
