# systemd/udp2sock@.socket
[Unit]
Description=Forward UDP port %i traffic to Unix socket

[Socket]
ListenDatagram=%i

[Install]
WantedBy=sockets.target
